# single-container/service/mariadb-maxscale-rtr/finish
#!/bin/sh
rm -f /var/lock/subsys/mariadb-maxscale-rtr