# single-container/service/mariadb-server-1/finish
#!/bin/sh
rm -f /var/lock/subsys/mariadb-server-1