# single-container/service/mariadb-server-1/run
#!/bin/sh
exec 2>&1
(mkdir -p /var/lock/subsys && touch /var/lock/subsys/mariadb-server-1)
exec /usr/local/mysql/bin/mysqld --defaults-file=/etc/mariadb-server-1.cnf
