# single-container/service/mariadb-maxscale-cdc/run
#!/bin/sh
exec 2>&1
(mkdir -p /var/lock/subsys && touch /var/lock/subsys/mariadb-maxscale-cdc)
exec /usr/bin/maxscale --nodaemon --user=maxscale --config=/etc/mariadb-maxscale-cdc.cnf